# deploy/systemd/cortex.service
[Unit]
Description=Cortex Agent Orchestrator
After=openclaw-gateway.service

[Service]
Type=simple
ExecStart=%h/projects/cortex/cortex --config %h/projects/cortex/cortex.toml
Restart=always
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=default.target
